# debian/systemd/opensips-restarter.path
[Path]
PathModified=/opt/irontec/itsbc/config/opensips/opensips_custom_listeners.cfg

[Install]
WantedBy=multi-user.target
