# debian/systemd/rtpengine-restarter.path
[Path]
PathModified=/opt/irontec/itsbc/config/rtpengine

[Install]
WantedBy=multi-user.target
